# packages/i2pd-git/i2pd.service
[Unit]
Description=I2P Router written in C++
Wants=network.target
After=network.target

[Service]
Type=simple
User=i2p
Group=i2p
EnvironmentFile=/etc/default/i2pd.conf
ExecStart=/usr/bin/i2p "$I2P_ARGS"
KillSignal=SIGINT
Restart=always

[Install]
WantedBy=multi-user.target
